# meta-mct/meta-s8050/recipes-phosphor/watchdog/phosphor-watchdog/timeout@.service
[Unit]
Description=Phosphor Watchdog

[Service]
Type=oneshot
ExecStart=/usr/bin/timeout %I

[Install]
WantedBy=basic.target