# meta-mct/meta-common/recipes-amd/debug-tool/amd-debug-tool/yaapd.service
[Unit]
Description=YAAP for AMD HDT

[Service]
Type=exec
ExecStart=/usr/bin/yaapd
SyslogIdentifier=yaap

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
